# 1_imageHub/supervisor/supervisord.conf
[supervisord]
nodaemon=true

[program:selenium-hub]
command=/root/selenium-hub.sh
