# 1_imageHub/supervisor/supervisord.conf
[supervisord]
pidfile=/tmp/supervisord.pid
logfile=/tmp/supervisord.log
nodaemon=true

[program:selenium-hub]
user=developer
command=/opt/selenium/selenium-hub.sh
