# 1_imageHub/supervisor/supervisord.conf
[supervisord]
nodaemon=true

[program:selenium-hub]
user=developer
command=/opt/selenium/selenium-hub.sh
