# 3_node-firefox/supervisor/supervisord.conf
[supervisord]
nodaemon=true

[program:xvfb]
command=/opt/xvfb.sh

[program:selenium-firefox-node]
command=/root/selenium-firefox-node.sh

