# 3_node-firefox/supervisor/supervisord.conf
[supervisord]
nodaemon=true

[program:xvfb]
user=developer
command=/opt/selenium/xvfb.sh

[program:selenium-firefox-node]
user=developer
command=/opt/selenium/selenium-firefox-node.sh

