# 4_node-chrome/supervisor/supervisord.conf
[supervisord]
nodaemon=true

[program:xvfb]
command=/opt/xvfb.sh

[program:selenium-chrome-node]
command=/root/selenium-chrome-node.sh

