# 4_node-chrome/supervisor/supervisord.conf
[supervisord]
nodaemon=true

[program:xvfb]
user=developer
command=/opt/selenium/xvfb.sh

[program:selenium-chrome-node]
#user=developer
command=/opt/selenium/selenium-chrome-node.sh

